# projects/RPi/devices/RPi4-GPICase/packages/dock/udev.d/99-hotdock.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/hotdock"